# files/usr/lib/systemd/system/heliumos-addons.service
[Unit]
Description=HeliumOS Addons mount service

[Service]
Type=oneshot
ExecStart=/usr/bin/make -C /usr/addons refresh

[Install]
WantedBy=multi-user.target
