# files/usr/lib/systemd/system/timedate-config.service
[Unit]
Description=HeliumOS timedate configuration service

[Service]
Type=oneshot
ExecStart=/usr/bin/timedatectl set-local-rtc 0
ExecStart=/usr/bin/timedatectl set-ntp 1

[Install]
WantedBy=multi-user.target
