# 10/system_root/etc/profile.d/timedate.sh
timedatectl set-ntp 1
timedatectl set-local-rtc 0